# etc/local.d/cpu_governor.start
#!/bin/bash

# set cpu governors to "ondemand"
for c in $(ls -d /sys/devices/system/cpu/cpu[0-9]*); do
	echo "ondemand" > $c/cpufreq/scaling_governor;
done

